# bin/etc/services.d/nginx/run
#!/command/execlineb -P
with-contenv
s6-envuidgid nginx
/usr/sbin/nginx
