# bin/etc/services.d/datareporter/run
#!/command/execlineb -P
with-contenv
s6-envuidgid redash
cd /app
./bin/docker-entrypoint server
