# bin/etc/services.d/plywood/run
#!/command/execlineb -P
with-contenv
s6-envuidgid redash
node /app/plywood/dist/server.js
