# bin/etc/services.d/plywood/run
#!/command/execlineb -P
with-contenv
s6-envuidgid redash
cd /app/plywood/
node dist/server.js
